# bat/install/rc.local
#!/bin/sh -e
#
# rc.local

/bat/irswatchdog

exit 0
